# deploy/docker/supervisord.conf
; ─────────────────────────────────────────────────────────────────────────────
;  supervisord.conf — process manager for the queue/scheduler container
;
;  ── Chosen approach ──────────────────────────────────────────────────────────
;  We run php-fpm in its OWN container (Dockerfile.api CMD ["php-fpm"]) and use
;  this supervisord config in a SEPARATE "worker" container built from the same
;  image. That keeps web request serving isolated from background processing, so:
;    • A flood of jobs can't starve web requests (and vice-versa).
;    • Workers can be scaled independently (`docker compose up --scale worker=N`).
;    • A crashing worker never takes down php-fpm.
;
;  This supervisord supervises:
;    • queue:work        — N Laravel queue workers (auto-restarted).
;    • schedule:work     — the Laravel scheduler (long-running, replaces cron;
;                          internally runs `schedule:run` every minute).
;
;  To use this file, the worker service overrides the image command to:
;      supervisord -c /etc/supervisor/conf.d/supervisord.conf
;  (mount this file there). See docker-compose.yml `worker` service.
; ─────────────────────────────────────────────────────────────────────────────

[supervisord]
nodaemon=true                 ; run in the foreground as PID 1's child
user=www-data
logfile=/dev/null             ; supervisor's own log -> discard
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

; ── Queue workers ────────────────────────────────────────────────────────────
; Processes jobs: MatchRequirementsToProperty, queued mail & notifications.
;   --sleep=3        wait 3s when the queue is empty before polling again
;   --tries=3        retry a failing job up to 3 times, then -> failed_jobs
;   --max-time=3600  recycle the worker hourly to bound memory growth
;   --max-jobs=1000  also recycle after 1000 jobs
;   --backoff=5      wait 5s before retrying a released job
; numprocs=2 runs two workers in THIS container; scale the container for more.
[program:queue-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --max-time=3600 --max-jobs=1000 --backoff=5
directory=/var/www/html
user=www-data
numprocs=2
autostart=true
autorestart=true
stopwaitsecs=3600             ; give an in-flight job up to an hour to finish
                              ; (Laravel stops cleanly on SIGTERM between jobs)
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; ── Scheduler ────────────────────────────────────────────────────────────────
; `schedule:work` is the containerized replacement for a system cron entry of
; `* * * * * php artisan schedule:run`. It stays running and fires due tasks
; each minute. Exactly ONE instance must run cluster-wide.
[program:scheduler]
command=php /var/www/html/artisan schedule:work
directory=/var/www/html
user=www-data
numprocs=1
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
